# resources/launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.reclaude.indexer</string>
  <key>ProgramArguments</key>
  <array>
    <string>__ELECTRON_PATH__</string>
    <string>__APP_PATH__</string>
    <string>--reindex</string>
    <string>--no-window</string>
  </array>
  <key>StartInterval</key>
  <integer>10800</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/indexer.out.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/indexer.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RECLAUDE_HEADLESS</key>
    <string>1</string>
  </dict>
</dict>
</plist>
